# rocky/nextcloud.service
[Unit]
Description=Nextcloud stack (Podman) — clean start/stop with prune
After=network-online.target podman.service
Wants=network-online.target
Requires=podman.service

[Service]
Type=oneshot
RemainAfterExit=yes
# Set project name so labels & network names are predictable
Environment=COMPOSE_PROJECT_NAME=nextcloud
Environment=COMPOSE_FILE=/opt/nextcloud/docker-compose.yml
WorkingDirectory=/opt/nextcloud

# --- PRE-START CLEANUP (handles unclean previous shutdowns) ---
# Try a graceful compose down first
ExecStartPre=/bin/sh -lc '/usr/bin/podman-compose -f "$COMPOSE_FILE" down --timeout 25 --remove-orphans || true'
# Force-remove any lingering containers for this project (labels are set by podman-compose)
ExecStartPre=/bin/sh -lc '/usr/bin/podman ps -aq --filter "label=io.podman.compose.project=$COMPOSE_PROJECT_NAME" | xargs -r /usr/bin/podman rm -f || true'
# Remove the project network if it still exists (common if down aborted)
ExecStartPre=/bin/sh -lc '/usr/bin/podman network ls --format "{{.Name}}" | grep -x "${COMPOSE_PROJECT_NAME}_web" >/dev/null 2>&1 && /usr/bin/podman network rm "${COMPOSE_PROJECT_NAME}_web" || true'
# Prune unused stuff (dangling images/containers/networks)
ExecStartPre=/usr/bin/podman system prune -f

# --- START ---
ExecStart=/usr/bin/podman-compose -f "$COMPOSE_FILE" up -d

# --- STOP + POST-STOP CLEANUP ---
ExecStop=/usr/bin/podman-compose -f "$COMPOSE_FILE" down --timeout 25 --remove-orphans
ExecStopPost=/usr/bin/podman system prune -f

# If compose up fails (non-zero), allow restart
Restart=on-failure
TimeoutStartSec=0
TimeoutStopSec=90

[Install]
WantedBy=multi-user.target