# debian/nextcloud.service
[Unit]
Description=Nextcloud stack (Docker) — clean start/stop with prune
After=network-online.target docker.service
Wants=network-online.target
Requires=docker.service

[Service]
Type=oneshot
RemainAfterExit=yes
# Set project name so labels & network names are predictable
Environment=COMPOSE_PROJECT_NAME=nextcloud
Environment=COMPOSE_FILE=/opt/nextcloud/docker-compose.yml
WorkingDirectory=/opt/nextcloud

# --- PRE-START CLEANUP (handles unclean previous shutdowns) ---
# Try a graceful compose down first
ExecStartPre=/bin/sh -lc '/usr/bin/docker-compose -f "$COMPOSE_FILE" down --timeout 25 --remove-orphans || true'
# Force-remove any lingering containers for this project
ExecStartPre=/bin/sh -lc '/usr/bin/docker ps -aq --filter "label=com.docker.compose.project=$COMPOSE_PROJECT_NAME" | xargs -r /usr/bin/docker rm -f || true'
# Remove the project network if it still exists
ExecStartPre=/bin/sh -lc '/usr/bin/docker network ls --format "{{.Name}}" | grep -x "${COMPOSE_PROJECT_NAME}_web" >/dev/null 2>&1 && /usr/bin/docker network rm "${COMPOSE_PROJECT_NAME}_web" || true'

# --- START ---
ExecStart=/usr/bin/docker-compose -f "$COMPOSE_FILE" up -d

# --- STOP ---
ExecStop=/bin/sh -lc '/usr/bin/docker-compose -f "$COMPOSE_FILE" down --timeout 25 --remove-orphans || true'

[Install]
WantedBy=multi-user.target
